# debian/postinst
#!/bin/bash -e

echo "Enabling MALOS service"
systemctl enable matrix-malos-zwave

echo "Please restart your Raspberry Pi after installation"
