# debian/postinst
#!/bin/bash -e

echo "Enabling MALOS Zwave service"
systemctl enable matrixio-malos-zwave

echo "Please restart your Raspberry Pi after installation"
